# dist/systemd/system/systemd-timesyncd.service.d/wait.conf
[Unit]
Wants=systemd-timesyncd-wait.socket systemd-timesyncd-wait.service
After=systemd-timesyncd-wait.socket

[Service]
Sockets=
ExecStart=
ExecStart=!!/lib/systemd/systemd-timesyncd-wrap /lib/systemd/systemd-timesyncd
ExecStopPost=/bin/rm -f -- /run/timesyncd/notify.sock
ExecStartPre=/bin/rm -f -- /run/timesyncd/notify.sock
ExecStartPre=+/bin/chown systemd-timesync /run/timesyncd
ReadWritePaths=/run/timesyncd
NotifyAccess=all
